# configure/files/etc/udev/rules.d/50-wol.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="enp*s0", RUN+="/usr/bin/ethtool -s %k wol g"
